# systemd/battery-idle-shutdown.timer
[Unit]
Description=Poll idle+battery state every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=30s

[Install]
WantedBy=timers.target
