# systemd/battery-idle-shutdown.service
[Unit]
Description=Check idle+battery state and power off if thresholds are met

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery-idle-shutdown.sh
