# systemd/swayidle.service
[Unit]
Description=Idle tracker for battery-idle-shutdown
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/swayidle-wrapper.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
